# src/lucid_component_projector/systemd/lucid-projector-helper.service
[Unit]
Description=LUCID projector helper (root-owned, owns RS232 serial port)
After=network-online.target
Before=lucid-agent-core.service
PartOf=lucid-agent-core.service

[Service]
Type=simple
User=root
Group=root
RuntimeDirectory=lucid
RuntimeDirectoryMode=0755
ExecStart=/home/lucid/lucid-agent-core/venv/bin/lucid-projector-helper
Restart=always
RestartSec=5
TimeoutStartSec=15
SyslogIdentifier=lucid-projector-helper

[Install]
WantedBy=multi-user.target
